# inittab
# LEVELS
si::sysinit:/sbin/rc sysinit
rc::bootwait:/sbin/rc boot
0:0:wait:/sbin/rc shutdown 
1:1:wait:/sbin/rc single
3:3:wait:/sbin/rc default
6:6:wait:/sbin/rc reboot

# TERMINALS
tty1:13:respawn:/sbin/agetty -J 38400 tty1 linux
