# init.d/ntpclient
#!/sbin/runscript
: ${NTP_HOST:=pool.ntp.org}
: ${NTP_FREQVAL:=0}

command="/usr/bin/ntpclient"
command_background=yes
command_args="-f ${NTP_FREQVAL} -s -l -i 600 -g 10000 -h ${NTP_HOST}"
command_user=nobody
command_group=nogroup
pidfile="/run/ntpclient.pid"

start_pre() {
	setcap CAP_SYS_TIME=eip "${command}"
}
